# co2.service
[Unit]
Description=Start CO2 Sensor

[Service]
ExecStart=sudo /home/pi/co2-sensor/start.sh

[Install]
WantedBy=multi-user.target